# deploy/Typora.ini
[program:Typora]
directory = /home/www/Typora
command = /home/www/Typora/venv/bin/python3 main.py --env=prod --port=8101
autorestart = true
stdout_logfile = /var/log/supervisor/Typora.log
stderr_logfile = /var/log/supervisor/Typora.err